# config/autostart/bitcoin.desktop
[Desktop Entry]
Type=Application
Name=Bitcoin
Exec=bash -c "ln -s /live/persistence/TailsData_unlocked/Persistent/.bitcoin /home/amnesia; /live/persistence/TailsData_unlocked/dotfiles/.local/bin/bitcoin-qt-wrapper -min -chain=main -startupnotify=chainstate-preload"
Terminal=false
Hidden=false
